# runtimes/litespeed-php/configs/supervisor/bagisto.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:lsphp]
command=/usr/sbin/gosu www-data /usr/local/lsws/lsphp83/bin/lsphp -b 127.0.0.1:9000
environment=PHP_LSAPI_CHILDREN="35",LSAPI_CHILDREN="35",LSAPI_AVOID_FORK="200M"
autostart=true
autorestart=true
startsecs=3
priority=5
stdout_logfile=/var/log/supervisor/lsphp.log
stderr_logfile=/var/log/supervisor/lsphp-error.log

[program:openlitespeed]
command=/bin/bash -c "/usr/local/lsws/bin/lswsctrl start; exec tail -F /usr/local/lsws/logs/error.log /usr/local/lsws/logs/access.log"
autostart=true
autorestart=true
startsecs=5
priority=10
stdout_logfile=/var/log/supervisor/openlitespeed.log
stderr_logfile=/var/log/supervisor/openlitespeed-error.log
